# home/.config/autostart/aw-qt.desktop
[Desktop Entry]
Name=ActivityWatch
Exec=aw-qt
Type=Application
X-GNOME-Autostart-enabled=true

